# server/deploy/systemd/pocket-claude.service
[Unit]
# Pocket Claude server — main process (FastAPI + uvicorn).
# Installed by deploy/install-linux.sh and run under the service user
# `pocket-claude`. The tunnel is a separate unit
# (pocket-claude-tunnel.service OR cloudflared.service depending on setup).
Description=Pocket Claude — Personal Claude Chat Server
Documentation=https://github.com/joshtech90/PocketClaude
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=@SERVICE_USER@
Group=@SERVICE_USER@
WorkingDirectory=@PROJECT_DIR@
EnvironmentFile=@PROJECT_DIR@/.env
# Important: the `claude` binary must be in the service user's PATH. We
# explicitly add /usr/local/bin (default npm-global prefix) in case systemd
# starts with a sparse PATH.
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
ExecStart=@PROJECT_DIR@/.venv/bin/python -m pocket_claude

# Robustness
Restart=always
RestartSec=5
StartLimitInterval=600
StartLimitBurst=10

# When `claude login` hasn't been run yet, the first request crashes with
# CLINotFoundError or similar. We do NOT want to spin in an infinite restart
# loop — respect the burst limiter; it pauses after 10 crashes in 10 min.

StandardOutput=journal
StandardError=journal
SyslogIdentifier=pocket-claude

# Hardening — defense in depth.
#
# NOTE: ProtectHome must stay OFF. The Claude CLI reads its OAuth credentials
# from /home/<service-user>/.claude/.credentials.json. With ProtectHome=true,
# that directory is masked for the service even when ReadWritePaths is also
# set — ReadWritePaths only widens write permissions, it does NOT undo the
# home-directory masking. Symptom: the service starts cleanly, the CLI spawns,
# but the SDK's control_request:initialize handshake silently times out after
# 60 s with no useful error. Confirmed empirically 2026-05-19.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ReadWritePaths=@PROJECT_DIR@
ReadWritePaths=/home/@SERVICE_USER@

[Install]
WantedBy=multi-user.target
