# server/deploy/systemd/pocket-claude.service
[Unit]
# Pocket Claude Server — Hauptprozess (FastAPI + uvicorn).
# Wird via deploy/install-linux.sh installiert und unter dem Service-User
# `pocket-claude` ausgeführt. Tunnel läuft als separater Service
# (pocket-claude-tunnel.service ODER cloudflared.service je nach Setup).
Description=Pocket Claude — Personal Claude Chat Server
Documentation=https://github.com/joshtech90/PocketClaude-Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=@SERVICE_USER@
Group=@SERVICE_USER@
WorkingDirectory=@PROJECT_DIR@
EnvironmentFile=@PROJECT_DIR@/.env
# Wichtig: `claude`-Binary muss im PATH des Service-Users liegen. Wir setzen
# /usr/local/bin (default npm-global-prefix) explizit drauf, falls systemd
# den PATH knapp hält.
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
ExecStart=@PROJECT_DIR@/.venv/bin/python -m pocket_claude

# Robustheit
Restart=always
RestartSec=5
StartLimitInterval=600
StartLimitBurst=10

# Wenn `claude login` noch fehlt, crasht der Server beim ersten Request mit
# CLINotFoundError o.ä. Wir wollen aber NICHT in eine endlose Restart-Schleife
# rennen, also den Burst-Limiter respektieren — der lässt nach 10 Crashs in
# 10 Min Pause.

StandardOutput=journal
StandardError=journal
SyslogIdentifier=pocket-claude

# Hardening — defense in depth
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ReadWritePaths=@PROJECT_DIR@
# Service-User-Home für ~/.claude/credentials.json beschreibbar lassen
ReadWritePaths=/home/@SERVICE_USER@

[Install]
WantedBy=multi-user.target
