# server/scripts/pocket-claude.service
[Unit]
Description=Pocket Claude — Personal Claude Chat Server
After=network.target

[Service]
Type=simple
User=@SERVICE_USER@
WorkingDirectory=@PROJECT_DIR@
EnvironmentFile=@PROJECT_DIR@/.env
ExecStart=@PROJECT_DIR@/.venv/bin/python -m pocket_claude
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=@PROJECT_DIR@

[Install]
WantedBy=multi-user.target
